# templates/mdbook_build.cron.j2
0  * * * * root /usr/local/sbin/mdbook-build >/dev/null
30 * * * * root /usr/local/sbin/mdbook-build >/dev/null
